# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:php7.2-fpm]
command=/usr/sbin/php-fpm7.2 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
